# roles/jenkins-pf-master/test/integration/default/serverspec/files/expected_jenkins_cleanup.timer
[Unit]
Description=Timer to cleanup the Jenkins workspace and old build
[Timer]
OnCalendar=Mon-Fri *-*-* 07:00:00
Unit=jenkins_cleanup.service
[Install]
WantedBy=multi-user.target
